# systemd/blipblop-node.service
[Unit]
Description=Blipblop Node Agent
Documentation=https://github.com/amimof/blipblop
After=network-online.target containerd.service blipblop-server.service
Wants=network-online.target containerd.service

[Service]
Type=simple

# Node reads its identity from /etc/blipblop/node.yaml by default.
# It connects to the server over gRPC on port 5743 with TLS.
ExecStart=/usr/local/bin/blipblop-node \
  --port=5743 \
  --tls-ca=/etc/blipblop/tls/ca.crt \
  --metrics-host=0.0.0.0

Restart=on-failure
RestartSec=5s

# Hardening (tune as needed)
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
