# systemd/blipblop-server.service
[Unit]
Description=Blipblop Server
Documentation=https://github.com/amimof/blipblop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Adjust binary path if installed elsewhere
ExecStart=/usr/local/bin/blipblop-server \
  --log-level=info \
  --tls-certificate=/etc/blipblop/tls/server.crt \
  --tls-key=/etc/blipblop/tls/server.key \
  --tls-ca=/etc/blipblop/tls/ca.crt \
  --tls-host=0.0.0.0 \
  --tcp-tls-host=0.0.0.0 \
  --metrics-host=0.0.0.0

Restart=on-failure
RestartSec=5s

# Hardening (tune as needed)
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
